# .build/DEBIAN/prerm
#!/bin/sh
systemctl stop homeseer
systemctl disable homeseer
